# rootfs/etc/services.d/httpd/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting httpd..."

exec python3 -m http.server 8020 --directory /config/public_html