# rootfs/etc/services.d/weewxd/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting weewxd..."

exec /home/weewx/bin/weewxd --config=/config/weewx.conf